# crontab
SHELL=/bin/bash
50 * * * * bash /app/email.bash &> test.txt
