# crontab
SHELL=/bin/bash
00 * * * * bash /app/email.bash &> test.txt
#* * * * * echo "LOL" >> /app/test.txt
